# .local/runit/sv/libinput-gestures/run
#!/bin/sh

[ -r ./conf ] && . ./conf

if [ $XDG_SESSION_TYPE=wayland ]; then
	exec /usr/bin/libinput-gestures ${OPTS:-} 2>&1
else
	env SVDIR=$HOME/.local/runit/run sv down libinput-gestures
fi
