# .local/runit/sv/calibre-catwatch/run
#!/usr/bin/env bash

[ -r ./conf ] && . ./conf


catalogCreate() {
  # Only run `calibredb catalog` when Calibre is closed
  while true; do
    if pgrep -fl calibre-parallel >/dev/null; then
      # Add `sleep` to reduce cpu usage
      echo ":: Waiting for running Calibre instances to close..."
      sleep 5 &&
      true
    else
      # Give time for `calibre` to properly close
      sleep 5 &&
      echo ":: Generating BibTeX file..."
      calibredb catalog "${CALIBREBIBTEX}" --entry-type=mixed --fields="${BIBTEXFIELDS}" -v &&
      # Remove first 2 lines in CALIBREBIBTEX as it was reserved for `@preamble`
      echo ":: Removing autogenerated '@preamble' from BibTeX file..."
      sed -i '1,2d' "${CALIBREBIBTEX}";
      echo ":: BibTeX file generated at ${CALIBREBIBTEX}"
      break
    fi
  done
}

caLibWatch() {
    while true; do 
      inotifywait -e modify "${CALIBREDB}" &&
      catalogCreate
    done
}

cd "${HOME}" &&
  caLibWatch "${OPTS:-}" 2>&1
