# .local/runit/sv/waydroid-session/finish
#!/bin/sh
set -e
exec /usr/bin/waydroid session stop
